# btrfs-snapshot.service
[Unit]
Description = Manage btrfs snapshots

[Service]
Type = oneshot
ExecStart = /usr/sbin/btrfs-snapshot create

[Install]
WantedBy = multi-user.target
