# package/usr/lib/systemd/system/btrfs-snapshots-manager.timer
[Unit]
Description = Manage btrfs snapshots

[Timer]
OnCalendar = daily
Persistent = true
WakeSystem = true
RandomizedDelaySec = 1h

[Install]
WantedBy = timers.target
