# package/usr/lib/systemd/system/btrfs-snapshot.service
[Unit]
Description = Manage btrfs snapshots

[Service]
Type = oneshot
ExecStart = /usr/bin/btrfs-snapshot create

[Install]
WantedBy = multi-user.target
