# package/usr/lib/systemd/system/btrfs-snapshots-manager.service
[Unit]
Description = Manage btrfs snapshots

[Service]
Type = oneshot
ExecStart = /usr/bin/btrfs-snapshots-manager

[Install]
WantedBy = multi-user.target
